# diamond.service
# This service is installed in systemd at /etc/systemd/system/diamond.service
# IMPORTANT: If diamond.py is renamed, update ExecStart below and reinstall:
#   sudo cp diamond.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl restart diamond.service

[Unit]
Description=Diamond Robot Controller
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/diamond
Environment="PATH=/home/admin/diamond/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/admin/diamond/.venv/bin/python /home/admin/diamond/diamond.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
