# diamond.service
# This service is installed in systemd at /etc/systemd/system/diamond.service
# Install or refresh after changing the drive command:
#   sudo cp diamond.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl restart diamond.service

[Unit]
Description=Diamond Rover Xbox Drive
Wants=bluetooth.service
After=bluetooth.service network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/diamond
Environment="PATH=/home/admin/diamond/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/admin/diamond/.venv/bin/python /home/admin/diamond/diamond.py --max-speed 1.00
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
