# diamond.service
[Unit]
Description=Diamond Robot Controller
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/diamond
Environment="PATH=/home/admin/diamond/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/admin/diamond/.venv/bin/python /home/admin/diamond/diamond.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
